# files/system/usr/lib/systemd/user/xwayland-satellite.service.d/10-custom.conf
[Unit]
Wants=xsettingsd.service

[Service]
Type=simple
TimeoutStartSec=300
